# crontab
32 * * * * /home/jbearer/Documents/auto-waiter/run_daemon.py
